# startup/alerting.timer
[Unit]
Description=run wailin alerting daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
